# service/log/run
#!/bin/sh
exec envdir ./env sh -c 'exec setuidgid ${UIDGID} multilog t ${NUM} ${SIZE} ./main'
